# deploy/3drisk-mp.service
[Unit]
Description=3D Risk multiplayer server (authoritative WebSocket host)
After=network.target

[Service]
Type=simple
User=iain
Environment=PORT=8787
Environment=NODE_ENV=production
# Optional: override the reconnect grace window (ms); default 5 min.
# Environment=MP_RECONNECT_MS=300000
WorkingDirectory=/srv/3drisk-mp
ExecStart=/usr/bin/env node /srv/3drisk-mp/server.cjs
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
